# app/crontab
5 * * * * cd /app && python app.py >> /var/log/cron.log 2>&1
